# vm.resources/files/ip_change_mon.service
[Unit]
Description=IP Change Monitor
Wants=network.target
After=network-online.target

[Service]
ExecStart=:/bin/bash -c "ip mon addr | sed -nu -r
\'s/.*[[:digit:]]+:[[:space:]]+([^[:space:]]+).*/\\1/p\' | while read iface; do
/resources/scripts/ensure_default_route.sh; done"

[Install]
WantedBy=multi-user.target default.target